# docker/supervisord.conf
[supervisord]
nodaemon=true
logfile=/dev/stdout
logfile_maxbytes=0
pidfile=/tmp/supervisord.pid

# Expose a control socket so `supervisorctl status|restart|stop|start <name>`
# works inside the container (useful for ops/debugging without restarting
# the whole container).
[unix_http_server]
file=/var/run/supervisor.sock
chmod=0700

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///var/run/supervisor.sock

[program:web]
directory=/app/backend
command=/opt/venv-sigma/bin/python -m uvicorn app.main:app --host 127.0.0.1 --port 8000
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
environment=PATH="/opt/venv-sigma/bin:/usr/local/texlive/current/bin/x86_64-linux:%(ENV_PATH)s",SIGMA_USERDATA_DIR="/app/userdata",SIGMA_JUPYTER_BIN="/opt/venv-jupyter/bin/jupyter",SIGMA_TEXLIVE_YEAR="%(ENV_SIGMA_TEXLIVE_YEAR)s",SIGMA_TEXLIVE_ROOT="/usr/local/texlive",SIGMA_TLMGR_BIN="/usr/local/texlive/current/bin/x86_64-linux/tlmgr",SIGMA_UPDATE_TLMGR_BIN="/usr/local/texlive/current/update-tlmgr-latest.sh"

[program:worker]
directory=/app/backend
command=/opt/venv-sigma/bin/huey_consumer app.workers.huey_tasks.huey -w 4 -k thread
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
environment=PATH="/opt/venv-sigma/bin:/usr/local/texlive/current/bin/x86_64-linux:%(ENV_PATH)s",SIGMA_USERDATA_DIR="/app/userdata",SIGMA_JUPYTER_BIN="/opt/venv-jupyter/bin/jupyter",SIGMA_TEXLIVE_YEAR="%(ENV_SIGMA_TEXLIVE_YEAR)s",SIGMA_TEXLIVE_ROOT="/usr/local/texlive",SIGMA_TLMGR_BIN="/usr/local/texlive/current/bin/x86_64-linux/tlmgr",SIGMA_UPDATE_TLMGR_BIN="/usr/local/texlive/current/update-tlmgr-latest.sh"

[program:nginx]
command=/usr/sbin/nginx -g "daemon off;"
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
